# macos-observer/launchd/com.carthing.btlink.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.carthing.btlink</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/exrector/Documents/ПРОЕКТЫ/carthing-release-integration/macos-observer/.build/debug/CarThingBTAudioCap</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>CARTHING_BT_HEADLESS</key>
        <string>1</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>/private/tmp/carthing-btlink.err</string>
</dict>
</plist>
