# macos-observer/launchd/com.carthing.btlink.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.carthing.btlink</string>
    <key>ProgramArguments</key>
    <array>
        <string>__CARTHING_BTLINK_EXECUTABLE__</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>CARTHING_BT_HEADLESS</key>
        <string>1</string>
        <key>CARTHING_BT_TCP_PORT</key>
        <string>49500</string>
        <key>CARTHING_ASSISTANT_PYTHON</key>
        <string>__CARTHING_ASSISTANT_PYTHON__</string>
        <key>CARTHING_ASSISTANT_SCRIPT</key>
        <string>__CARTHING_ASSISTANT_SCRIPT__</string>
        <key>VA_DENOISE</key>
        <string>0</string>
        <key>VA_STT_BACKEND</key>
        <string>mlx</string>
        <key>VA_STT_LANGUAGE</key>
        <string>auto</string>
        <key>VA_STT_MLX_MODEL</key>
        <string>__CARTHING_ASSISTANT_MODEL__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>StandardOutPath</key>
    <string>/private/tmp/carthing-btlink.out</string>
    <key>StandardErrorPath</key>
    <string>/private/tmp/carthing-btlink.err</string>
</dict>
</plist>
