# reference/legacy-mfi-iap2/device_root/etc/udev/rules.d/97-mfi.rules
SUBSYSTEM=="module" DEVPATH=="/module/apple_mfi_auth_i2c" ACTION=="add" RUN+="/bin/add-mfi-dev"
SUBSYSTEM=="module" DEVPATH=="/module/apple_mfi_auth_i2c" ACTION=="remove" RUN+="/bin/remove-mfi-dev"
