# experiment/mfi-chip/reference/device_root/etc/supervisor.d/iap2_agent.conf
[program:iap2_agent]
command=/home/superbird/slot_a/iap2_agent
environment=IAP2_TEST_HID_USAGE="0x00E9",IAP2_TEST_HID_MODE="4"
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/tmp/iap2.log
